# etc/skel/.config/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Install LimaLinux
Exec=sudo calamares
Icon=calamares
Terminal=false
Categories=Qt;System;
